# scripts/macos/com.wakeupnpc.check-new-memes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.wakeupnpc.check-new-memes</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>node</string>
      <string>__REPO_PATH__/scripts/check-new-memes.js</string>
      <string>--notify</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_PATH__</string>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__REPO_PATH__/logs/check-new-memes.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_PATH__/logs/check-new-memes.err.log</string>
  </dict>

</plist>
